# daemon/xpum.service.in
[Unit]
Description=XPUM daemon

After=syslog.target

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=/usr/local/xpum:/$LD_LIBRARY_PATH"
ExecStart=@CPACK_PACKAGING_INSTALL_PREFIX@/bin/xpumd  --pid_file /var/xpum_daemon.pid
ExecStop=/bin/kill -TERM $MAINPID
User=root

[Install]
WantedBy=multi-user.target