# rest/xpum_rest.service.in
[Unit]
Description=XPUM_Rest daemon

After=xpum.target

[Service]
Type=oneshot
RemainAfterExit=yes

RuntimeDirectory=xpum

ExecStart=@CPACK_PACKAGING_INSTALL_PREFIX@/rest/xpum_rest.sh start
ExecStop=@CPACK_PACKAGING_INSTALL_PREFIX@/rest/xpum_rest.sh stop
ExecReload=@CPACK_PACKAGING_INSTALL_PREFIX@/rest/xpum_rest.sh reload

User=xpum

[Install]
WantedBy=multi-user.target
